# opsight-agent/opsight-agent.service
[Unit]
Description=Opsight Server Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/opsight-agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
